# systemd/service-a.service
[Unit]
Description=Service A - Public Gateway
After=network.target service-b.service service-c.service
Requires=service-b.service service-c.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/production-services/services/service-a
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=3
Environment=NODE_ENV=production
Environment=SERVICE_B_URL=http://service-b.internal:3002
StandardOutput=journal
StandardError=journal
SyslogIdentifier=service-a

[Install]
WantedBy=multi-user.target
