# systemd/service-c.service
[Unit]
Description=Service C - Internal Processor
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/production-services/services/service-c
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=3
Environment=NODE_ENV=production
Environment=SERVICE_A_URL=http://service-a.internal:3001
StandardOutput=journal
StandardError=journal
SyslogIdentifier=service-c

[Install]
WantedBy=multi-user.target
